# dsmr_reader/rootfs/etc/services.d/addon/run
#!/usr/bin/with-contenv bashio

while true; do
  ADDON_VERSON=$(bashio::addon.version)
  bashio::log.green "Home Assistant DSMR Reader Add-on - Release: ${ADDON_VERSON}"
  sleep 2
done
